# crontab
*/10 * * * * run-one rq -all --queue default